# ft_services/srcs/mysql/srcs/supervisord.conf
[supervisord]
nodaemon=true

[program:mysql]
command=sh setup.sh
autostart=true
autorestart=true